# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn hello:app -b localhost:5000
directory=/deploy/app
